# build/etc/supervisord.conf
[supervisord]
user=root
nodaemon=true

[program:php-fpm]
command=/usr/sbin/php-fpm7 --nodaemonize

[program:nginx]
command=/bin/rutorrent.sh
